# rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
# set default error log and run in foreground
nginx -e "/var/run/s6/nginx-error-log-fifo" -g "daemon off;"
