# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh
# use with-contenv to pass env vars to phpfpm
php-fpm8 --nodaemonize
